# homework/homelab/deployment/postgres/pg-backup.timer
[Unit]
Description=Portable Postgres Backup

[Timer]
OnCalendar=3d
AccuracySec=1h 

[Install]
WantedBy=timers.target

